# samba/rootfs/etc/s6-overlay/s6-rc.d/smbd/run
#!/command/with-contenv bash

exec smbd ${SMBD_ARGS:--F --no-process-group}
