# samba/rootfs/etc/s6-overlay/s6-rc.d/nmbd/run
#!/command/with-contenv bash

exec nmbd ${NMBD_ARGS:--F --no-process-group}
